# bin/storm_service
#!/bin/sh
### BEGIN INIT INFO
# Provides:          storm services
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Description:       storm services
# Short-Description: Enable storm services including hdfs
### END INIT INFO

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
STORM_BIN=$STORM_HOME/bin
NAME=storm
DESC=storm
USER=adragomi
ROTATE_SUFFIX=
test -x $STORM_BIN || exit 0
RETVAL=0
set -e
cd /

start_storm () {
    set +e
    $STORM_BIN/storm nimbus &
    sleep 5
    $STORM_BIN/storm supervisor &
    sleep 5
    $STORM_BIN/storm ui &
    case "$?" in
      0)
        echo SUCCESS
        RETVAL=0
        ;;
      1)
        echo TIMEOUT - check /var/log/storm/startup_log
        RETVAL=1
        ;;
      *)
        echo FAILED - check /var/log/storm/startup_log
        RETVAL=1
        ;;
    esac
    set -e
}

stop_storm () {
    set +e
    if [ $RETVAL = 0 ] ; then
        ps axwwww | grep "[s]torm.daemon.nimbus" | awk '{print $1}' | xargs kill -9
        ps axwwww | grep "[s]torm.daemon.supervisor" | awk '{print $1}' | xargs kill -9
        ps axwwww | grep "[s]torm.ui.core" | awk '{print $1}' | xargs kill -9
        RETVAL=$?
        if [ $RETVAL != 0 ] ; then
            echo FAILED - check /var/log/storm/shutdown_log
        fi
    else
        echo No nodes running
        RETVAL=0
    fi
    set -e
}

restart_storm() {
    stop_storm
    start_storm
}

case "$1" in
    start)
        echo -n "Starting $DESC: "
        start_storm
        echo "$NAME."
        ;;
    stop)
        echo -n "Stopping $DESC: "
        stop_storm
        echo "$NAME."
        ;;
    force-reload|restart)
        echo -n "Restarting $DESC: "
        restart_storm
        echo "$NAME."
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|force-reload}" >&2
        RETVAL=1
        ;;
esac
exit $RETVAL

